# deploy/systemd/reverb.service
[Unit]
Description=Laravel Reverb Websocket Server
After=network.target
Wants=network-online.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/onlyverified
Environment=APP_ENV=production
Environment=APP_DEBUG=false
Restart=always
RestartSec=5
ExecStart=/usr/bin/php /var/www/onlyverified/artisan reverb:start --host=127.0.0.1 --port=8080 --no-interaction
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

